# deploy/golive/contest-start.service
[Unit]
Description=BNB Hack Track1 — automated contest-start (panic -> reset -> live)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/contest-start.sh
