# deploy/golive/golive-remind@.service
[Unit]
Description=BNB Hack go-live Telegram reminder (%i)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/golive-remind.sh %i
